# profile/airootfs/etc/xdg/autostart/kernos-installer-online.desktop
[Desktop Entry]
Type=Application
Name=Install KernOS
Exec=kernos-calamares-online
Icon=calamares
Terminal=false
X-GNOME-Autostart-enabled=true
